# dls_sst.cron
00 00 * * * /syndata/staging/dls-prod/gfs_sel_url_generator.sh
00 05 * * * /syndata/staging/dls-prod/dls-dlscript/gfs_fcst_sel_aria_downloader.sh 00
30 04 * * * /syndata/staging/dls-prod/dls-dlscript/visibility_fcst_downloader.sh 00
30 06 * * * /syndata/staging/dls-prod/dls-dlscript/ww3_fcst_downloader.sh 00
30 10 * * * /syndata/staging/dls-prod/dls-dlscript/visibility_fcst_downloader.sh 06
00 11 * * * /syndata/staging/dls-prod/dls-dlscript/gfs_fcst_sel_aria_downloader.sh 06
30 12 * * * /syndata/staging/dls-prod/dls-dlscript/ww3_fcst_downloader.sh 06
30 16 * * * /syndata/staging/dls-prod/dls-dlscript/visibility_fcst_downloader.sh 12
45 16 * * * /syndata/staging/dls-prod/dls-dlscript/hycom_fcst_downloader.sh uv 123
05 17 * * * /syndata/staging/dls-prod/dls-dlscript/hycom_fcst_downloader.sh sst 123
00 17 * * * /syndata/staging/dls-prod/dls-dlscript/gfs_fcst_sel_aria_downloader.sh 12
00 18 * * * /syndata/staging/dls-prod/dls-dlscript/hycom_fcst_downloader.sh uv 4
15 18 * * * /syndata/staging/dls-prod/dls-dlscript/hycom_fcst_downloader.sh sst 4
30 18 * * * /syndata/staging/dls-prod/dls-dlscript/ww3_fcst_downloader.sh 12
30 22 * * * /syndata/staging/dls-prod/dls-dlscript/visibility_fcst_downloader.sh 18
35 22 * * * /syndata/staging/dls-prod/dls-dlscript/hycom_fcst_downloader.sh uv 5
45 22 * * * /syndata/staging/dls-prod/dls-dlscript/hycom_fcst_downloader.sh sst 5
00 23 * * * /syndata/staging/dls-prod/dls-dlscript/gfs_fcst_sel_aria_downloader.sh 18
30 00 * * * /syndata/staging/dls-prod/dls-dlscript/ww3_fcst_downloader.sh 18
