# mdls.cron
PYTHONIOENCODING=utf8
00 00 * * * /syndata/dls-prod/gfs_sel_url_generator.sh
30 04 * * * /syndata/dls-prod/dls-dlscript/visibility_fcst_downloader.sh 00
00 05 * * * /syndata/dls-prod/dls-dlscript/gfs_fcst_sel_aria_downloader.sh 00
30 05 * * * /syndata/dls-prod/dls-dlscript/generate_land_reports.sh 00
30 06 * * * /syndata/dls-prod/dls-dlscript/ww3_fcst_downloader.sh 00
45 06 * * * /syndata/dls-prod/dls-dlscript/generate_reports.sh 00
15 10 * * * /syndata/dls-prod/dls-dlscript/ss_fcst_downloader.sh 0
30 10 * * * /syndata/dls-prod/dls-dlscript/visibility_fcst_downloader.sh 06
00 11 * * * /syndata/dls-prod/dls-dlscript/gfs_fcst_sel_aria_downloader.sh 06
30 11 * * * /syndata/dls-prod/dls-dlscript/generate_land_reports.sh 06
30 12 * * * /syndata/dls-prod/dls-dlscript/ww3_fcst_downloader.sh 06
45 12 * * * /syndata/dls-prod/dls-dlscript/generate_reports.sh 06
04 16 * * 0 /syndata/dls-prod/dls-dlscript/generate_wsr.sh
15 16 * * * /syndata/dls-prod/dls-dlscript/ss_fcst_downloader.sh 1
30 16 * * * /syndata/dls-prod/dls-dlscript/visibility_fcst_downloader.sh 12
00 17 * * * /syndata/dls-prod/dls-dlscript/gfs_fcst_sel_aria_downloader.sh 12
30 17 * * * /syndata/dls-prod/dls-dlscript/generate_land_reports.sh 12
30 18 * * * /syndata/dls-prod/dls-dlscript/ww3_fcst_downloader.sh 12
45 18 * * * /syndata/dls-prod/dls-dlscript/generate_reports.sh 12
00 19 * * * /syndata/dls-prod/cleanup_sst.sh
30 22 * * * /syndata/dls-prod/dls-dlscript/visibility_fcst_downloader.sh 18
00 23 * * * /syndata/dls-prod/dls-dlscript/gfs_fcst_sel_aria_downloader.sh 18
30 23 * * * /syndata/dls-prod/dls-dlscript/generate_land_reports.sh 18
30 00 * * * /syndata/dls-prod/dls-dlscript/ww3_fcst_downloader.sh 18
45 00 * * * /syndata/dls-prod/dls-dlscript/generate_reports.sh 18
00 01 * * * /syndata/dls-prod/cleanup_gfs.sh
00 01 * * * /syndata/dls-prod/cleanup_ww3.sh
00 01 * * * /syndata/dls-prod/cleanup_reports.sh
00 01 * * * /syndata/dls-prod/tcleanup_reports.sh
*/30 * * * *  /syndata/dls-prod/dls-dlscript/send_scheduled_reports.sh
*/30 * * * *  /syndata/dls-prod/dls-dlscript/srsend_scheduled_reports.sh
45 09,11,21,23 * * * /syndata/dls-prod/dls-dlscript/msend_scheduled_reports.sh
45 09,11,21,23 * * * /syndata/dls-prod/dls-dlscript/srmsend_scheduled_reports.sh
15 03,15 * * * /syndata/dls-prod/dls-dlscript/msend_scheduled_reports.sh
15 03,15 * * * /syndata/dls-prod/dls-dlscript/srmsend_scheduled_reports.sh
30 03,15 * * * /syndata/dls-prod/dls-dlscript/msend_scheduled_reports.sh
30 03,15 * * * /syndata/dls-prod/dls-dlscript/srmsend_scheduled_reports.sh
30 09,21 * * * /syndata/dls-prod/dls-dlscript/gmo_msend_scheduled_reports.sh
30 09,21 * * * /syndata/dls-prod/dls-dlscript/srgmo_msend_scheduled_reports.sh
30 00,04,06,10,12,16,18,22 * * * /syndata/dls-prod/dls-dlscript/gmo_msend_scheduled_reports.sh
30 00,04,06,10,12,16,18,22 * * * /syndata/dls-prod/dls-dlscript/srgmo_msend_scheduled_reports.sh
15 05,11,17,23 * * * /syndata/dls-prod/dls-dlscript/gmo_msend_scheduled_reports.sh
15 05,11,17,23 * * * /syndata/dls-prod/dls-dlscript/srgmo_msend_scheduled_reports.sh
45 03 * * * /syndata/dls-prod/dls-dlscript/gmo_msend_scheduled_reports.sh
45 03 * * * /syndata/dls-prod/dls-dlscript/srgmo_msend_scheduled_reports.sh
00 05,11,17,23 * * * /syndata/dls-prod/dls-dlscript/gmo_msend_scheduled_reports.sh
00 05,11,17,23 * * * /syndata/dls-prod/dls-dlscript/srgmo_msend_scheduled_reports.sh
00 * * * * /syndata/dls-prod/dls-dlscript/check_ta.sh
30 * * * * /syndata/dls-prod/dls-dlscript/check_ta.sh
00 22 * * 0 /syndata/dls-prod/dls-dlscript/send_wsr.sh
